# etc/cron.d/ufw_backup
0 2 * * * root /usr/local/bin/ufw_backup.sh
